# src/main/python/dac-bridge.service
[Unit]
Description=DAC Bridge - MQTT to DAC8532 for HRV control
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/dac-bridge --mqtt-host localhost --mqtt-port 1883
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Hardening
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
