# systemd/dac-bridge.service
[Unit]
Description=HRV Bridge - MQTT to PWM/DAC for HRV control
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
# GPIO 5: bypass valve (digital), GPIO 12+13: PWM outputs (HW PWM)
# GPIO 17, 18, 22, 23 are reserved for Waveshare AD/DA board
ExecStart=/usr/local/bin/hrv-bridge --mqtt-host localhost --gpio17 5 --gpio18 12 --gpio19 13 --verbose
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
