# systemd/dac-bridge.service
[Unit]
Description=HRV Bridge - MQTT to PWM/DAC for HRV control
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
# PWM on GPIO 18 and GPIO 19
ExecStart=/usr/local/bin/hrv-bridge --mqtt-host localhost --gpio18 18 --gpio19 19 --verbose
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
